# md-check@.service
[Unit]
Description=Scrub RAID array %i

[Service]
ExecStart=/usr/sbin/mdscrub %i
